# crontab
* * * * * docker container exec moodle-webserver php /var/www/html/admin/cli/cron.php >/dev/null